# systemd/backup@.timer
[Unit]
Description=Timer unit for %i backup.

[Timer]
OnCalendar=hourly
