# systemd/backup.service
[Unit]
Description=Make backups
EnvironmentFile=-/tmp/agent-info

[Service]
Type=oneshot
ExecStart=/usr/bin/backup
