# systemd/backup.timer
[Unit]
Description=Timer unit for backup service.

[Timer]
OnCalendar=hourly
